# gatherClear.cron
#!/bin/sh

cd ~/ib/
./gather cron
./clear
